# debian/DEBIAN/postinst
#!/bin/bash

echo "Install x3 wifi drivers"

systemctl enable hobot-wifi.service

exit 0

