# debian/DEBIAN/postinst
#!/bin/bash

systemctl enable hobot-wifi.service
systemctl enable hobot-bluetooth.service

exit 0
