# debian/DEBIAN/postinst
#!/bin/bash

systemctl enable hobot-wifi.service

exit 0
