# debian/DEBIAN/postinst
#!/bin/bash

echo "Install x3 sdb wifi drivers"

depmod

systemctl enable hobot-wifi.service

exit 0

